# etc/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:elasticsearch]
directory=/data
command=/elasticsearch/bin/elasticsearch -Des.config=/data/elasticsearch.yml
autostart=yes
autorestart=yes